# .launchd/com.user.free-yangmao.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- 任务标识符 -->
  <key>Label</key>
  <string>com.user.free-yangmao.update</string>

  <!-- 要执行的脚本 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/wenjiaqi/Downloads/free-yangmao/update.sh</string>
  </array>

  <!-- 执行时间：每天上午9点 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- 标准输出重定向到日志文件 -->
  <key>StandardOutPath</key>
  <string>/Users/wenjiaqi/Downloads/free-yangmao/logs/launchd-out.log</string>

  <!-- 错误输出重定向到日志文件 -->
  <key>StandardErrorPath</key>
  <string>/Users/wenjiaqi/Downloads/free-yangmao/logs/launchd-err.log</string>

  <!-- 工作目录 -->
  <key>WorkingDirectory</key>
  <string>/Users/wenjiaqi/Downloads/free-yangmao</string>

  <!-- 环境变量 -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- 如果上一次任务还在运行，跳过本次执行 -->
  <key>AbandonProcessGroup</key>
  <true/>

</dict>
</plist>
